# docs/service/systemd/ambient-swarm@.service
[Unit]
Description=Ambient Swarm (watch) for %i
After=network.target

[Service]
Type=simple
WorkingDirectory=%i
Environment=PYTHONUNBUFFERED=1
# Optional kill-switches:
# Environment=AMBIENT_PAUSED=0
# Environment=AMBIENT_MAX_PROPOSALS_PER_HOUR=0
ExecStart=/usr/local/bin/ambient watch %i --approval-mode webhook
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target

